# docker/cron/crontabs
# minute hour day_of_month month day_of_week command_to_run
*/15 8,12 * * 1-5 /bin/sh /problem-car-generator >/dev/null 2>&1